# deploy/ferg-api.service
[Unit]
Description=Ferg Design API
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/ferg
EnvironmentFile=/opt/ferg/server/.env
ExecStart=/usr/bin/npm --prefix server run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
